# macos/com.josh.cloudsessioncleaner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.josh.cloudsessioncleaner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>if [ $(du -sk ~/Library/Application\ Support/CloudDocs/session 2&gt;/dev/null | awk '{print $$1}') -gt 5000000 ] &amp;&amp; rm -rf ~/Library/Application\ Support/CloudDocs/session &amp;&amp; osascript -e 'display notification "Deleted iCloud session cache &gt;5GB" with title "Auto Cleanup"'; fi</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer> <!-- Monday -->
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
